# Resources/Packaging/postinstall
#!/bin/zsh --no-rcs
#
# postinstall
#
# Bundled into the App Auto-Patch installer pkg's Scripts payload and executed
# as root by installer(8)/MDM once the package's scripts archive is expanded.
#
# App Auto-Patch already knows how to self-install: when the script detects
# it's running from anywhere other than its canonical install folder
# (/Library/Management/AppAutoPatch) or the /usr/local/bin CLI symlink, it
# copies itself into place, (re)creates the LaunchDaemon, and hands off the
# actual discovery/patch run to that LaunchDaemon via restart_aap().
#
# This script simply stages the bundled copy of the script to a temporary
# location and runs it once, so that existing self-install logic takes over.
# No changes to the App Auto-Patch script itself are required.
#
# The staged script is run synchronously (foreground, no backgrounding/
# disown) so the installer/MDM transaction does not report success until
# App Auto-Patch has actually finished copying itself into place and setting
# up its LaunchDaemon. This step is normally quick - the staged instance
# calls restart_aap(), which itself backgrounds a brief, disowned
# `launchctl bootstrap` and exits immediately, so it hands off the real
# discovery/patch run to the LaunchDaemon rather than running it here.
#
# Backgrounding+disowning this outer invocation was tried previously so the
# installer transaction could return immediately, but installer(8)/MDM can
# tear down this script's process group as soon as postinstall itself
# returns - killing the still-running staged script (and racing our own
# cleanup below) before it finished copying itself into place. Running it
# synchronously avoids that: nothing below runs, and the temp folder is not
# removed, until the staged script has actually exited.

set -e

SCRIPT_DIR="$(cd "$(dirname "$0")" && pwd)"
TMP_DIR=$(mktemp -d /private/tmp/AppAutoPatchInstaller.XXXXXX)

cp "${SCRIPT_DIR}/App-Auto-Patch-via-Dialog.zsh" "${TMP_DIR}/App-Auto-Patch-via-Dialog.zsh"
chmod 755 "${TMP_DIR}/App-Auto-Patch-via-Dialog.zsh"

"${TMP_DIR}/App-Auto-Patch-via-Dialog.zsh"

rm -rf "${TMP_DIR}"

exit 0
